# supervisord.conf
[supervice]
loglevel=INFO

[program:test]
command=python3 test_loop.py
autostart=true
autorestart=true
startsecs=2
stdout_logfile=test.log
stderr_logfile=test_err.log
